# src/main/debian/prerm
#!/usr/bin/env bash

# Stop running instance if needed
systemctl stop gridify.service

# Disable service if exists
systemctl disable gridify.service

# remove symlink
rm /usr/bin/gridifyd
